# infra/low-end-ci/systemd/storycraft-ci-prune.service
[Unit]
Description=Docker system prune (StoryCraft CI)

[Service]
Type=oneshot
ExecStart=/usr/bin/docker system prune -f
