# infra/low-end-ci/systemd/storycraft-forgejo.service
[Unit]
Description=StoryCraft Forgejo (eco-mode Git server)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=@USER@
Group=@USER@
Environment=HOME=@HOME@
Environment=FORGEJO_DATA_DIR=@HOME@/storycraft-ci/forgejo/data
Environment=FORGEJO_CONFIG_DIR=@HOME@/storycraft-ci/forgejo/config
WorkingDirectory=@HOME@/githubcursor/StoryCraft-Studio/infra/low-end-ci
ExecStart=/usr/bin/docker compose -f @HOME@/githubcursor/StoryCraft-Studio/infra/low-end-ci/docker-compose.forgejo.yml up -d
ExecStop=/usr/bin/docker compose -f @HOME@/githubcursor/StoryCraft-Studio/infra/low-end-ci/docker-compose.forgejo.yml down

[Install]
WantedBy=default.target
