# infra/low-end-ci/systemd/worldscript-ci-prune.timer
[Unit]
Description=Weekly Docker prune for WorldScript low-end CI

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
